# systemd/cya.service
[Unit]
Description=Run rsync

[Service]
User=root
Group=root
Type=oneshot
ExecStart=/home/USER/bin/cya save
