# systemd/cya.timer
[Unit]
Description=Run cya every week

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=multi-user.target
